# themes/luci-theme-material-mini/root/etc/uci-defaults/30_luci-theme-material-mini
#!/bin/sh

if [ "$PKG_UPGRADE" != 1 ]; then
	uci get luci.themes.MaterialMini >/dev/null 2>&1 || \
	uci batch <<-EOF
		set luci.themes.MaterialMini=/luci-static/material-mini
		set luci.main.mediaurlbase=/luci-static/material-mini
		commit luci
	EOF
fi

exit 0
