# themes/luci-theme-material-mini/root/etc/uci-defaults/30_luci-theme-material-mini
#!/bin/sh
uci batch <<-EOF
	set luci.themes.MaterialMini=/luci-static/material-mini
	set luci.main.mediaurlbase=/luci-static/material-mini
	commit luci
EOF
exit 0
